# templates/lobalancer@.timer.j2
[Unit]
Description=Run lobalancer every minute

[Timer]
OnCalendar=*:*:00

[Install]
WantedBy=timers.target